# .config/systemd/user/update-theme.path
[Unit]
Description=Triggers update-theme service when Xresources file changes

[Path]
PathChanged=%h/.Xresources
Unit=update-theme.service

[Install]
WantedBy=default.target
